# tsimcloud-backend.service
[Unit]
Description=TsimCloud Backend API Server
Documentation=https://github.com/tarikogut/tssystem
After=network.target mysql.service redis.service rabbitmq-server.service
Wants=mysql.service redis.service rabbitmq-server.service

[Service]
Type=simple
User=tsimcloud
Group=tsimcloud
WorkingDirectory=/opt/tsimcloud/backend
ExecStart=/opt/tsimcloud/backend/server
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tsimcloud-backend

# Environment variables
Environment=GIN_MODE=release
Environment=GO_ENV=production

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tsimcloud/backend/logs /opt/tsimcloud/backend/static

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096



[Install]
WantedBy=multi-user.target 